# systemd/coredns-fgpu.service
# CoreDNS for dynamic DNS synthesis (IP.<SUFFIX> → IP)
#
# CUSTOMIZE: <SUFFIX> is your DNS suffix (e.g., work, corp, lab)
# Filename should match: coredns-<SUFFIX>.service
# Requires: /etc/coredns/Corefile configured with your <SUFFIX>

[Unit]
Description=CoreDNS (dynamic .<SUFFIX> IP synthesizer)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/coredns -conf /etc/coredns/Corefile
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

