# systemd/expose-socks-tsFGPU.service
# Expose SOCKS5 proxy from namespace to primary tailnet
#
# CUSTOMIZE:
#   <NAMESPACE>: Your namespace name (e.g., tsWork, tsCorp)
#   <HOST_TAILNET_IP>: Your proxy host's Tailscale IP on primary tailnet
#   <SOCKS_PORT>: SOCKS5 port (default: 11080)
#   <VETH_NS_IP>: Namespace veth IP (default: 10.200.0.6)
# Filename should match: expose-socks-<NAMESPACE>.service

[Unit]
Description=Expose <NAMESPACE> SOCKS5 proxy on host tailnet IP
After=network-online.target ts-netns.service dante-<NAMESPACE>.service
Wants=network-online.target
Requires=ts-netns.service dante-<NAMESPACE>.service

[Service]
Restart=always
RestartSec=1
# Listen on primary tailnet IP and forward to namespace veth IP
ExecStart=/usr/bin/socat TCP-LISTEN:<SOCKS_PORT>,bind=<HOST_TAILNET_IP>,fork,reuseaddr TCP:<VETH_NS_IP>:<SOCKS_PORT>

[Install]
WantedBy=multi-user.target

