# systemd/socks-tsFGPU.service
[Unit]
Description=SOCKS5 proxy inside tsFGPU namespace
After=network-online.target tailscaled-tsFGPU.service
Wants=network-online.target tailscaled-tsFGPU.service
Requires=tailscaled-tsFGPU.service

[Service]
Restart=always
RestartSec=1
# Bind to the namespace's veth IP (reachable from host)
ExecStart=/usr/sbin/ip netns exec tsFGPU /usr/bin/microsocks -i 10.200.0.6 -p 11080

[Install]
WantedBy=multi-user.target

