# systemd/tailscaled-tsFGPU.service
[Unit]
Description=Tailscale daemon in tsFGPU netns (FarmGPU)
After=network-online.target ts-netns.service
Wants=network-online.target ts-netns.service
Requires=ts-netns.service

[Service]
Type=simple
ExecStartPre=/usr/sbin/ip netns exec tsFGPU ip link set lo up
ExecStart=/usr/sbin/ip netns exec tsFGPU /usr/sbin/tailscaled \
  --state=/var/lib/tailscale-tsFGPU/tailscaled.state \
  --socket=/run/tailscale-tsFGPU/tailscaled.sock \
  --port=41642
Restart=always

[Install]
WantedBy=multi-user.target

