# systemd/ts-netns.service
# Network namespace setup for isolated secondary tailnet
#
# Runs setup-ts-netns.sh which creates the namespace, veth pairs, and NAT rules.
# Edit /usr/local/sbin/setup-ts-netns.sh to configure namespace name and IPs.

[Unit]
Description=Setup network namespace and veth pairs for secondary tailnet
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/setup-ts-netns.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

